# profile/profile.d/aliases
# Make commands colorful
alias ls='ls -G -h'
alias less='less -R'
alias egrep='egrep --color=always --exclude-dir=.git --exclude-dir=.svn --exclude-dir=.hg --exclude-dir=.svn --binary-file=without-match'

alias f='find . -name'
alias df='df -h'
alias top='top -o cpu -s 1'
alias sedi='sed -i ""'

# Ensure we get the term we expect on the remote end
alias ssh='TERM=xterm-256color ssh'

# Get the weather forecast :D
alias forecast='curl wttr.in'

# Manage my dotfiles
alias dotfiles="git --git-dir=$HOME/.dotfiles/ --work-tree=$HOME"

# Update my dotfiles
alias updot="dotfiles --no-pager incoming && dotfiles rebase origin/master && dotfiles subup"
