# profile/profile.d/aliases
# Make commands colorful
alias ls='ls -G -h'
alias less='less -R'
alias aws='aws --color=on'

alias f='find . -name'         # simplified find for basic search in current directory
alias df='df -h'               # Human readable df by default
alias top='top -o cpu -s 1'    # top defaults to ordering by CPU usage and 1 second delay
alias sedi='sed -i ""'         # in-place sed

# Upgrade all Pip packages
alias pipupgrade="pip freeze --local | grep -v '^-e' | cut -d = -f 1 | sed '/vboxapi/d' | xargs pip install -U"

# Update my dotfiles
alias updot="git -C ${HOME}/dotfiles --no-pager incoming && git -C ${HOME}/dotfiles rebase origin/master && git -C ${HOME}/dotfiles subup"
