# profile/profile.d/aliases
# Make commands colorful
alias ls='ls -G -h'
alias less='less -R'
alias aws='aws --color=on'
alias egrep='egrep --color=always --exclude-dir=.git --exclude-dir=.svn --exclude-dir=.hg --exclude-dir=.svn --binary-file=without-match'

alias f='find . -name'
alias df='df -h'
alias top='top -o cpu -s 1'
alias sedi='sed -i ""'

alias cask='brew cask'

# Get the weather forecast :D
alias forecast='curl wttr.in'

# Upgrade all Pip packages
alias pipupgrade="pip freeze --local | grep -v '^-e' | cut -d = -f 1 | sed '/vboxapi/d' | xargs pip install -U"

# Update my dotfiles
alias updot="git -C ${HOME}/dotfiles --no-pager incoming && git -C ${HOME}/dotfiles rebase origin/master && git -C ${HOME}/dotfiles subup"
