# profile/profile.d/env
#!/bin/zsh
#----------------------------------------------------------------------
# Environment variable things
#----------------------------------------------------------------------
export EDITOR=vim
export GREP_OPTIONS='--color=always'
export PIP_DOWNLOAD_CACHE="$HOME/.pip/cache"
export SSH_ASKPASS='/usr/libexec/ssh-askpass'
export JAVA_HOME='/System/Library/Java/JavaVirtualMachines/1.6.0.jdk/Contents/Home'

#----------------------------------------------------------------------
# For go
#----------------------------------------------------------------------
export GOPATH="$HOME/code/go"
export PATH="$PATH:$GOPATH/bin"

#----------------------------------------------------------------------
# Load cabal (Haskell) functionality
#----------------------------------------------------------------------
[[ -e "$HOME/.cabal/bin" ]] && export PATH="$HOME/.cabal/bin:$PATH"
