# profile/profile.d/functions
# TODO make this work for bash also
myshell="$(ps -p $$ | awk 'NR > 1 { print $4 }')"
case "${myshell##*/}" in
  zsh)
    fpath+=( "${0%/*}/functions.d" )
    for file in $(find -s "${0%/*}/functions.d/" -depth 1 -type f -not -name '*.zwc'); do
      autoload "${file##*/}"
    done
    ;;
  bash)
    for file in $(find -s $(dirname $BASH_SOURCE)/functions.d -depth 1 -type f); do
      eval "function ${file##*/} { source \"${file}\"; }"
    done
    ;;
esac
