# profile/profile.d/functions
# TODO make this work for bash also

myshell="$(ps -p $$ | awk 'NR > 1 { print $4 }')"
case "$(basename -- $myshell)" in
  zsh)
    fpath+=( "$(dirname $0)/functions.d" )
    for file in "$(dirname $0)/functions.d/"*; do
      autoload "$(basename $file)"
    done
    ;;
  bash)
    for file in $(find $(dirname $BASH_SOURCE)/functions.d -depth 1 -type f); do
      eval "function $(basename $file) { source \"$file\"; }"
    done
    ;;
esac
